# crontab
0 * * * * rccw /opt/save/cleanup.sh                                                                                                  
*/5  *  *  *  *  rccw /opt/save/waterfall_2.sh 80m 0.3 17
