# crontab
0 * * * * rccw /opt/save/cleanup.sh                                                                                                  
*/5 * * * * rccw /opt/save/waterfall.sh 2m  
